# scripts/webserver/check_glusterfs.service
[Unit]
Description=Check and Reconfigure GlusterFS

[Service]
ExecStart=/usr/local/bin/configure_glusterfs.sh
